# packaging/linux/vbuff.service
[Unit]
Description=vbuff private clipboard manager
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'test -n "$WAYLAND_DISPLAY" || test -n "$DISPLAY"'
ExecStart=/usr/bin/vbuff
Restart=on-failure
RestartSec=2s
TimeoutStopSec=10s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=-%h/.local/share/vbuff -%h/.config/vbuff
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=graphical-session.target
